# roles/workstations-school-net/files/display-manager/gdm28/gdm.service.d/10-add-env.conf
[Service]
ExecStartPre=/usr/bin/extract-gnome-shell-resources
