# roles/common/files/lest/netshare-networld.automount
[Unit]
Description=Networld
ConditionPathExists=/netshare/networld
Requires=network-online.target
After=network-online.target

[Automount]
Where=/netshare/networld

[Install]
WantedBy=multi-user.target
