# roles/common/files/netshare-networld.automount
[Unit]
Description=Networld
ConditionPathExists=/netshare/networld

[Automount]
Where=/netshare/networld

[Install]
WantedBy=multi-user.target
