# roles/common/files/netshare-users.automount
[Unit]
Description=Users
ConditionPathExists=/netshare/users

[Automount]
Where=/netshare/users

[Install]
WantedBy=multi-user.target
